# scripts/macos/launchagents/com.story.goal.mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.story.goal.mcp</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>bash</string>
    <string>-lc</string>
    <string>exec story-goal-mcp</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>#{HOME}</string>
  <key>StandardOutPath</key>
  <string>${HOME}/Library/Logs/story-goal-mcp.out</string>
  <key>StandardErrorPath</key>
  <string>${HOME}/Library/Logs/story-goal-mcp.err</string>
</dict>
</plist>

